# conf.d/mk_xauth.conf
# for some reason it says "exit status 0; not expected"
# but it works fine
[program:mk_xauth]
command=/bin/bash -c "touch /root/.Xauthority && xauth add $DISPLAY /root $(mcookie)"